# templates/service/zookeeper_service.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          zookeeper services
# Required-Start:    $network
# Required-Stop:     $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Description:       zookeeper services
# Short-Description: Enable Zookeeper services
### END INIT INFO
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
ZOOKEEPER_BIN=<%= zookeeper_home %>/bin
NAME=zookeeper
DESC=zookeeper
USER=<%= user %>
ROTATE_SUFFIX=
test -x $ZOOKEEPER_BIN || exit 0
RETVAL=0
set -e
cd /

start_zookeeper () {
    set +e
    $ZOOKEEPER_BIN/zkServer.sh start
    case "$?" in
      0)
        echo SUCCESS
        RETVAL=0
        ;;
      1)
        echo TIMEOUT - check /var/log/zookeeper/zookeeper.log
        RETVAL=1
        ;;
      *)
        echo FAILED - check /var/log/zookeeper/zookeeper.log
        RETVAL=1
        ;;
    esac
    set -e
}

stop_zookeeper () {
    set +e
    if [ $RETVAL = 0 ] ; then
        $ZOOKEEPER_BIN/zkServer.sh stop
        RETVAL=$?
        if [ $RETVAL != 0 ] ; then
            echo FAILED - check /var/log/zookeeper/zookeeper.log
        fi
    else
        echo No nodes running
        RETVAL=0
    fi
    set -e
}

restart_zookeeper() {
    stop_zookeeper
    start_zookeeper
}

case "$1" in
    start)
        echo -n "Starting $DESC: "
        start_zookeeper
        echo "$NAME."
        ;;
    stop)
        echo -n "Stopping $DESC: "
        stop_zookeeper
        echo "$NAME."
        ;;
    force-reload|restart)
        echo -n "Restarting $DESC: "
        restart_zookeeper
        echo "$NAME."
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|force-reload}" >&2
        RETVAL=1
        ;;
esac
exit $RETVAL
