# templates/init.RedHat.erb
#!/bin/sh
#
### BEGIN INIT INFO
# Provides:          kafka
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts Kafka
# Description:       starts Kafka
### END INIT INFO

if [ `id -u` -ne 0 ]; then
echo "You need root privileges to run this script"
exit 1
fi

. /etc/init.d/functions

if [ -f /etc/sysconfig/kafka ]; then
. /etc/sysconfig/kafka
fi

# # The following variables can be overwritten in $DEFAULT
#
# Run kafka as this user ID and group ID
KAFKA_USER=kafka
KAFKA_GROUP=kafka

KAFKA_BASEDIR=/usr/local/kafka
KAFKA_CONFIG=/usr/local/kafka/config
KAFKA_LOG=/usr/local/kafka/logs
KAFKA_START=yes

JDK_DIRS="/usr/lib/jvm/default-java /usr/lib/jvm/java-6-sun /usr/lib/jvm/java-1.5.0-sun /usr/lib/j2sdk1.5-sun /usr/lib/j2sdk1.5-ibm /usr/lib/jvm/java-1.7.0-openjdk-amd64 /usr/lib/jvm/java-1.7.0-openjdk.x86_64"

# Look for the right JVM to use
for jdir in $JDK_DIRS; do
if [ -r "$jdir/bin/java" -a -z "${JAVA_HOME}" ]; then
JAVA_HOME="$jdir"
fi
done
export JAVA_HOME

# Default Java options
# Set java.awt.headless=true if JAVA_OPTS is not set so the
# Xalan XSL transformer can work without X11 display on JDK 1.4+
if [ -z "$JAVA_OPTS" ]; then
JAVA_OPTS="-Djava.awt.headless=true"
fi

# End of variables that can be overwritten in $DEFAULT

# overwrite settings from default file
if [ -f "$DEFAULT" ]; then
. "$DEFAULT"
fi

# If ulimit nofiles was set in $DEFAULT,
# then change the limit for this process now.
# This will not set nofiles higher than
# system limit.
if [ -n "$KAFKA_NOFILES_ULIMIT" ]; then
ulimit -n "$KAFKA_NOFILES_ULIMIT"
fi

# Define other required variables
KAFKA_PID="/var/run/$NAME.pid"
DAEMON="${JAVA_HOME}/bin/java"

CMD="$KAFKA_BASEDIR/bin/kafka-run-class.sh kafka.Kafka $KAFKA_CONFIG/server.properties >> $KAFKA_LOG/server.log 2>&1"

kafka_sh() {
    # Escape any double quotes in the value of JAVA_OPTS
    JAVA_OPTS="$(echo $JAVA_OPTS | sed 's/\"/\\\"/g')"
    KAFKA_OPTS="${KAFKA_OPTS} ${JAVA_OPTS}"

    # Export Kafka environment variables from $DEFAULT
    export KAFKA_CONFIG JMX_PORT KAFKA_JMX_OPTS SCALA_VERSION KAFKA_JVM_PERFORMANCE_OPTS KAFKA_HEAP_OPTS KAFKA_LOG4J_OPTS KAFKA_OPTS

    daemon --user=$KAFKA_USER --pidfile=$KAFKA_PID "$CMD"
    status="$?"
    return $status
}

start() {
    if [ -z "$JAVA_HOME" ]; then
        echo -n $"no JDK found - please set JAVA_HOME "
        exit 1
    fi

    if [ -n "$KAFKA_START" -a "$KAFKA_START" != "yes" ]; then
        echo -n $"KAFKA_START not set to 'yes' in $DEFAULT, not starting "
        exit 0
    fi

    echo -n $"Starting $prog: "
    kafka_sh

    RETVAL=$?
    [ $RETVAL -eq 0 ] && touch /var/lock/subsys/kafka && success || failure
    echo
    return $RETVAL
}

stop() {
    echo -n $"Shutting down $prog: "
    killproc -p $KAFKA_PID kafka
    RETVAL=$?
    echo
    [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/kafka
    return $RETVAL
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart|reload|force-reload)
    stop
    start
  ;;
condrestart)
  [ -f /var/lock/subsys/kafka ] && restart || :
  ;;
status)
  status -p $PIDFILE kafka
  RETVAL=$?
  ;;
*)
  echo "Usage: $0 {start|stop|status|restart|reload|force-reload|condrestart}"
  RETVAL=1
esac