# templates/zookeeper.service.erb
# Autogenerated by Puppet from zookeeper/templates/zookeeper.service.erb
# Do not edit

[Service]
ExecStart=/bin/sh -c 'set -x; . /etc/zookeeper/conf/environment; CLASSPATH="/usr/lib/zookeeper/zookeeper.jar:/usr/lib/zookeeper/lib/*:$CLASSPATH"; CLASSPATH="$(. /usr/lib/zookeeper/bin/zkEnv.sh ; echo $CLASSPATH)"; mkdir -p /var/log/zookeeper; $JAVA "-Dzookeeper.log.dir=/var/log/zookeeper" "-Dzookeeper.root.logger=INFO,CONSOLE" -cp "$CLASSPATH" $JVMFLAGS "$ZOOMAIN" "$ZOOCFG"'
Restart=always
SyslogIdentifier=zookeeper
User=<%= @user %>
Group=<%= @group %>

[Install]
WantedBy=multi-user.target
